# systemd/remote-studio.service
[Unit]
Description=Remote Studio Background Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/remote-studio/daemon/remote_studio_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
